# systemd/drop-ins/sonar-monitor.service.d/50-fd-watchdog-override.conf
# EMSN Sonar FD-watchdog override - applied 2026-05-05
#
# De algemene 10-emsn-baseline.conf zet LimitNOFILE=4096, TasksMax=128.
# Voor sonar-monitor en sonar-bavaria willen we juist de strenge originele
# waarden behouden: een FD-leak (zoals het paho/socket leak incident van
# 2026-04-22, 3 dagen 0 detecties) moet binnen uren zichtbaar worden als
# crash + restart, niet als stille degradatie.
#
# Steady-state verbruik: ~7 FDs / 7 threads. LimitNOFILE=512 geeft ruime
# marge maar crasht binnen uren bij een lek.
[Service]
LimitNOFILE=512
TasksMax=50
