# systemd/bat-monitor.service
[Unit]
Description=EMSN Bats Monitor - Vleermuisdetectie met BatDetect2
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=ronny
WorkingDirectory=/home/ronny/emsn-bats
ExecStart=/home/ronny/emsn-bats/venv/bin/python3 -m scripts.detection.bat_monitor
Restart=on-failure
RestartSec=10
WatchdogSec=300

Environment=PYTHONUNBUFFERED=1
Environment=NUMBA_DISABLE_JIT=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bat-monitor

[Install]
WantedBy=multi-user.target
