# systemd/bat-sync.service
[Unit]
Description=EMSN Bats PostgreSQL Sync - SQLite naar PostgreSQL
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn-bats
ExecStart=/home/ronny/emsn-bats/venv/bin/python3 -m scripts.sync.pg_sync

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bat-sync
