# systemd/batty-sync.service
[Unit]
Description=EMSN Bats BattyBirdNET-Bavaria PostgreSQL Sync
After=network-online.target batty-watcher.service
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn-bats
ExecStart=/home/ronny/emsn-bats/venv/bin/python3 -m scripts.sync.batty_pg_sync

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=batty-sync
