# systemd/emsn-bats-web.service
[Unit]
Description=EMSN Bats Web UI - Vleermuismonitoring Dashboard
After=network.target

[Service]
Type=simple
User=ronny
WorkingDirectory=/home/ronny/emsn-bats
ExecStart=/home/ronny/emsn-bats/venv/bin/python3 -m web.app
Restart=on-failure
RestartSec=10

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=emsn-bats-web

[Install]
WantedBy=multi-user.target
