# systemd/emsn-sonar-web.service
[Unit]
Description=EMSN Sonar Web UI - Vleermuismonitoring Dashboard
After=network.target

[Service]
Type=simple
User=ronny
WorkingDirectory=/home/ronny/emsn-sonar
ExecStart=/home/ronny/emsn-sonar/venv/bin/python3 -m web.app
Restart=on-failure
RestartSec=10

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=emsn-sonar-web

[Install]
WantedBy=multi-user.target
