# systemd/sonar-cleanup.service
[Unit]
Description=EMSN Sonar Cleanup - Verwijder oude audio en spectrogrammen

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn-sonar
ExecStart=/home/ronny/emsn-sonar/venv/bin/python3 -m scripts.archive.cleanup

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sonar-cleanup

# Resource guards - oneshot dat veel kleine file deletes doet.
LimitNOFILE=256
TasksMax=20
MemoryMax=128M
