# systemd/sonar-ha-discovery.service
[Unit]
Description=EMSN Sonar Home Assistant MQTT Discovery
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn-sonar
ExecStart=/home/ronny/emsn-sonar/venv/bin/python3 -m scripts.monitoring.ha_mqtt_discovery

Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sonar-ha-discovery

[Install]
WantedBy=multi-user.target
