# systemd/sonar-batdetect2-sync.timer
[Unit]
Description=EMSN Sonar PostgreSQL Sync Timer - Elk uur op :15

[Timer]
OnCalendar=*-*-* *:15:00
Persistent=true

[Install]
WantedBy=timers.target
