# systemd/sonar-monitor.service
[Unit]
Description=EMSN Sonar Monitor - Vleermuisdetectie met BatDetect2
After=network.target sound.target
Wants=network.target

[Service]
Type=notify
NotifyAccess=main
User=ronny
WorkingDirectory=/home/ronny/emsn-sonar
ExecStart=/home/ronny/emsn-sonar/venv/bin/python3 -m scripts.detection.sonar_monitor
Restart=on-failure
RestartSec=10
WatchdogSec=300
TimeoutStartSec=120

Environment=PYTHONUNBUFFERED=1
Environment=NUMBA_DISABLE_JIT=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sonar-monitor

# Resource guards. LimitNOFILE bewust laag (<<systeem default 1024)
# zodat een toekomstige FD-leak binnen uren ipv dagen zichtbaar wordt
# als crash + restart-counter, niet als stille degradatie. Werkelijk
# verbruik in steady state is ~7 FDs / 7 threads.
LimitNOFILE=512
TasksMax=50
MemoryMax=1500M

[Install]
WantedBy=multi-user.target
