# systemd/drop-ins/sonar-bavaria.service.d/50-fd-watchdog-override.conf
# EMSN Sonar FD-watchdog override - applied 2026-05-05
#
# Zelfde policy als sonar-monitor: 10-emsn-baseline.conf overschrijft
# de intentioneel lage limits. Deze drop-in herstelt ze (priority 50 > 10).
# Bavaria spawnt subprocess bat_ident.py per WAV - 50 tasks is ruim genoeg.
[Service]
LimitNOFILE=512
TasksMax=50
